# systemd/pi1-camera.service
[Unit]
Description=Pi5 Timelapse Camera Capture
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/2pis/pi5
ExecStart=/usr/bin/python3 /home/pi/2pis/pi5/camera.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
# Give the network time to fully come up before starting
ExecStartPre=/bin/sleep 5

[Install]
WantedBy=multi-user.target
