# systemd/pi3-camera.service
[Unit]
Description=Pi3 Timelapse Camera Capture
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/2pis/pi3
ExecStart=/usr/bin/python3 /home/pi/2pis/pi3/camera.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
# Give the network time to fully come up before starting
ExecStartPre=/bin/sleep 5

[Install]
WantedBy=multi-user.target
